# airootfs/etc/skel/.config/autostart/startx.desktop
[Desktop Entry]
Version=1.0
Name=start
Exec=core-brain
Type=Application
Terminal=false
Categories=Utility
